# base/etc/periodic/daily/spamassassin
#!/bin/sh
sa-update && kill -HUP $(cat /var/run/spamd.pid)
